# meta-co-pilot/recipes-co-pilot/diagnostics-service/diagnostics-service/diagnostics.service
[Unit]
Description=Co-Pilot diagnostics subsystem
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/co-pilot-diagnostics
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
