# meta-sigma-racer-wingman/recipes-sigma-racer-wingman/diagnostics-service/sigma-telemetry-relay/diagnostics.service
[Unit]
Description=Sigma Racer Wingman telemetry WiFi relay
Documentation=man:sigma-telemetry-relay(8)
After=network.target sigma-racer-vehicle.service
Wants=sigma-racer-vehicle.service

[Service]
Type=simple
EnvironmentFile=-/etc/sigma-racer-wingman/diagnostics.env
ExecStart=/usr/bin/sigma-telemetry-relay
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
