# meta-sigmaracer-wingman/recipes-sigmaracer-wingman/sigmaracer-wingman-services/sigmaracer-wingman-services/sigmaracer-wingman-ui.target
[Unit]
Description=Sigma Racer Wingman instrument cluster (sole UI)
Documentation=man:sigmaracer-wingman(7)
Requires=weston.service cluster-ui.service
After=weston.service
Before=navigation.service bluetooth.service camera.service

[Install]
WantedBy=graphical.target
