# meta-co-pilot/recipes-co-pilot/bluetooth-service/bluetooth-service/bluetooth.service
[Unit]
Description=Co-Pilot bluetooth subsystem
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/co-pilot-bluetooth
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
