# meta-co-pilot/recipes-co-pilot/vehicle-service/vehicle-service/vehicle.service
[Unit]
Description=Co-Pilot Vehicle Signal Abstraction (CAN / telemetry)
After=network.target
Before=cluster-ui.service
Wants=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/co-pilot/vehicle.env
ExecStart=/usr/bin/co-pilot-vehicle
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
