# meta-co-pilot/recipes-co-pilot/camera-service/camera-service/camera.service
[Unit]
Description=Co-Pilot camera subsystem
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/co-pilot-camera
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
