# meta-sigma-racer-wingman/recipes-sigma-racer-wingman/camera-service/camera-service/camera.service
[Unit]
Description=Sigma Racer Wingman camera subsystem
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sigma-racer-wingman-camera
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
