# meta-co-pilot/recipes-core/seatd/seatd/seatd.service
[Unit]
Description=Seat management daemon
Documentation=man:seatd(1)
After=systemd-user-sessions.service

[Service]
Type=simple
ExecStart=/usr/bin/seatd -g video
Restart=on-failure

[Install]
WantedBy=multi-user.target
