# meta-co-pilot/recipes-co-pilot/instrumentation/instrumentation/cluster-ui.service
[Unit]
Description=Sigma Instrument Cluster UI (instrumentation)
Documentation=file:/usr/share/doc/instrumentation/README
After=weston.service
Requires=weston.service
PartOf=co-pilot-ui.target
Conflicts=getty@tty7.service

[Service]
Type=simple
User=cluster
Group=cluster
EnvironmentFile=-/etc/co-pilot/ui.env
ExecStartPre=/bin/mkdir -p /run/user/cluster
ExecStartPre=/bin/chown cluster:cluster /run/user/cluster
ExecStart=/usr/bin/sigma-dash
Restart=always
RestartSec=1
WatchdogSec=30
StandardOutput=journal
StandardError=journal
StartLimitIntervalSec=60
StartLimitBurst=10

[Install]
WantedBy=co-pilot-ui.target
